# luci-app-passwall/root/etc/init.d/passwall_server
#!/bin/sh /etc/rc.common

START=99

boot_func() {
	config_load passwall_server
	config_get start_delay global start_delay 90
	case "$start_delay" in
		''|*[!0-9]*) start_delay=90 ;;
		*) [ "$start_delay" -le 600 ] 2>/dev/null || start_delay=90 ;;
	esac
	[ "$start_delay" -gt 0 ] && sleep "$start_delay"
	restart
}

boot() {
	boot_func >/dev/null 2>&1 &
}

start() {
	lua /usr/lib/lua/luci/passwall/server_app.lua start
}

stop() {
	lua /usr/lib/lua/luci/passwall/server_app.lua stop
}

restart() {
	stop
	start
}

reload() {
	echo "Warning: This service does not support configuration reloading"
	echo "Performing full restart instead..."
	restart
}
